# yunwifi/yunwifi-root/files/etc/hotplug.d/iface/98-yunwifi-wifidog
[ "$INTERFACE" = "wan" -a "$ACTION" = "ifup" ] && {
	pidof get_wifidog_conf.sh || /usr/sbin/get_wifidog_conf.sh &
}

