# yunwifi/yunwifi-root/files/etc/hotplug.d/iface/98-yunwifi-wifidog
[ "$INTERFACE" = "wan" ] && {
        /usr/sbin/get_wifidog_conf.sh &
        ipset -N WHITE nethash --hashsize 5000 --probes 2
        /usr/sbin/whitelist.sh update
}

