# yunwifi/yunwifi-root/files/etc/hotplug.d/iface/99-yunwifi-sysctl
[ "$INTERFACE" = "lan" ] && {
	iso=$(uci get wireless.@wifi-iface[0].isolate)
	if [ "$iso" = "1" ]; then
		sysctl -w net.ipv4.conf.$DEVICE.proxy_arp_pvlan=1
	fi
}

