# yunwifi/yunwifi-bwcontrol/files/etc/hotplug.d/iface/99-yunwifi-bw
[ "$INTERFACE" = "lan" -a "$ACTION" = "ifup" ] && {
	logger -t yunwifi-bwcontrol lan up
	/usr/sbin/redirectifb.sh $DEVICE
}
