# yunwifi/yunwifi-bwcontrol/files/etc/hotplug.d/iface/99-yunwifi-bw
[ "$INTERFACE" = "lan" ] && {
	logger "lan up"

	/usr/sbin/redirectifb.sh $DEVICE
}
