# yunwifi/wifidog-yunwifi/files/wifidog.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=65
EXTRA_COMMANDS="status"
EXTRA_HELP="        status Print the status of the service"

start() {
    mkdir -p /tmp/wifidog
    /usr/sbin/setifbroot.sh
	/usr/sbin/get_wifidog_conf.sh &
}

stop() {
	/usr/bin/wdctl stop
}

status() {
	/usr/bin/wdctl status
}
