# yunwifi/yunwifi-root/files/etc/init.d/hc5661-leds
#!/bin/sh /etc/rc.common
START=99
start(){
    /usr/sbin/hc5661_internet.sh &
}
stop() {
    logger "stop"
}
