# yunwifi/yunwifi-root/files/etc/init.d/hc5661-mac
#!/bin/sh /etc/rc.common
START=11
start(){
    /usr/sbin/hc5661_mac.sh
}
stop() {
    logger "stop"
}
