# yunwifi/yunwifi-root/files/etc/init.d/yunwifi-auth-interface
#!/bin/sh /etc/rc.common
START=15
start(){
    /usr/sbin/set_auth_interface.sh
}
stop() {
    logger "stop"
}
