# yunwifi/yunwifi-root/files/etc/init.d/yunwifi-whitelist
#!/bin/sh /etc/rc.common
START=99
start(){
    whitelist.sh do_table
	whitelist.sh update-with-dnsmasq
	/usr/sbin/check_dog.sh &
}
stop() {   
    logger "whitelist stop"
}
